# reports/master_governor/com.wanglifei.beijing-politics-master-governor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wanglifei.beijing-politics-master-governor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/wanglifei/GaokaoPolitics/beijing-politics-sync/scripts/run_master_governor_daily.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/wanglifei/GaokaoPolitics/beijing-politics-sync</string>

  <key>StandardOutPath</key>
  <string>/Users/wanglifei/Desktop/北京高考政治/reports/master_governor/launchd.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/wanglifei/Desktop/北京高考政治/reports/master_governor/launchd.stderr.log</string>
</dict>
</plist>

